# etc/sv/openggsn/run
#!/bin/sh
exec 2>&1
exec chrt 20 ggsn  -c /etc/osmocom/ggsn.conf -f -d
