# etc/sv/osmo-sgsn/run
#!/bin/sh
exec 2>&1
exec chrt 10 osmo-sgsn -c /etc/osmocom/osmo-sgsn.cfg
